# conf/etc/init/uwsgi.conf
# Emperor uWSGI script

description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [06]

setuid user
setgid www-data

respawn

exec uwsgi --emperor /etc/uwsgi
